# systemd/jackrabbit-dlm@.service
[Unit]
Description=JackrabbitDLM — Volatile Crypto Vault
Documentation=https://github.com/rapmd73/JackrabbitDLM
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/JackrabbitDLM
ExecStart=/usr/bin/python3 /home/JackrabbitDLM/JackrabbitDLM 0.0.0.0 37373
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/JackrabbitDLM/Logs /home/JackrabbitDLM/Disk /home/JackrabbitDLM/Quarantine
PrivateTmp=yes

# Resource limits
MemoryMax=256M
TasksMax=64

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jackrabbit-dlm

[Install]
WantedBy=multi-user.target
